# files/etc/uci-defaults/90-tollgate-nodogsplash-files
#!/bin/sh

echo "Running TollGate NoDogSplash file copy script..."

# Backup any existing nodogsplash files
mkdir -p /etc/nodogsplash/htdocs.bak
cp -f /etc/nodogsplash/htdocs/*.html /etc/nodogsplash/htdocs.bak/ 2>/dev/null || true
cp -f /etc/nodogsplash/htdocs/*.json /etc/nodogsplash/htdocs.bak/ 2>/dev/null || true

# Debug output
echo "Source files:"
ls -la /etc/tollgate/nodogsplash/htdocs/

# HTML files
if [ -f /etc/tollgate/nodogsplash/htdocs/splash.html ]; then
  echo "Copying splash.html..."
  cp -fv /etc/tollgate/nodogsplash/htdocs/splash.html /etc/nodogsplash/htdocs/
else
  echo "Warning: splash.html not found in source directory!"
fi

# JSON files
if [ -f /etc/tollgate/nodogsplash/htdocs/asset-manifest.json ]; then
  cp -fv /etc/tollgate/nodogsplash/htdocs/asset-manifest.json /etc/nodogsplash/htdocs/
fi
if [ -f /etc/tollgate/nodogsplash/htdocs/manifest.json ]; then
  cp -fv /etc/tollgate/nodogsplash/htdocs/manifest.json /etc/nodogsplash/htdocs/
fi

# Handle static files
if [ -d /etc/tollgate/nodogsplash/htdocs/static ]; then
  # Ensure destination directories exist
  mkdir -p /etc/nodogsplash/htdocs/static/css
  mkdir -p /etc/nodogsplash/htdocs/static/js
  mkdir -p /etc/nodogsplash/htdocs/static/media

  # CSS files
  if [ -d /etc/tollgate/nodogsplash/htdocs/static/css ]; then
    for css_file in $(find /etc/tollgate/nodogsplash/htdocs/static/css -type f 2>/dev/null); do
      cp -fv "$css_file" /etc/nodogsplash/htdocs/static/css/
    done
  fi

  # JS files
  if [ -d /etc/tollgate/nodogsplash/htdocs/static/js ]; then
    for js_file in $(find /etc/tollgate/nodogsplash/htdocs/static/js -type f 2>/dev/null); do
      cp -fv "$js_file" /etc/nodogsplash/htdocs/static/js/
    done
  fi

  # Media files
  if [ -d /etc/tollgate/nodogsplash/htdocs/static/media ]; then
    for media_file in $(find /etc/tollgate/nodogsplash/htdocs/static/media -type f 2>/dev/null); do
      cp -fv "$media_file" /etc/nodogsplash/htdocs/static/media/
    done
  fi
fi

# Restart nodogsplash to apply changes
/etc/init.d/nodogsplash restart

# Create a flag file to indicate we've completed this step
touch /etc/tollgate/nodogsplash-files-copied

exit 0