# files/etc/init.d/tollgate-basic
#!/bin/sh /etc/rc.common

START=99
STOP=15

EXTRA_COMMANDS="status"
EXTRA_HELP="        status          Check service status"

start() {
    echo "Starting TollGate basic" > /tmp/basic.log 2>&1
    tollgate-basic >> /tmp/basic.log 2>&1 &
    echo "PID: $!" >> /tmp/basic.log 2>&1
}

stop() {
    echo "Stopping TollGate basic"
    killall tollgate-basic
}

status() {
    if pgrep -f "tollgate-basic" >/dev/null; then
        echo "TollGate basic is running"
        return 0
    else
        echo "TollGate basic is not running"
        return 1
    fi
}
